# supervisor.conf
[program:web]
command=/usr/local/bin/gunicorn django.core.wsgi:get_wsgi_application --bind 0.0.0.0:8000
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
